# xos/synchronizer/supervisor/veg-observer.conf
[program:veg-observer]
command=python /opt/xos/synchronizers/veg/veg-synchronizer.py -C /opt/xos/synchronizers/veg/veg_synchronizer_config
